# backend/deploy/goinvoicer.service
[Unit]
Description=GoInvoicer API
After=network.target

[Service]
Type=simple
User=goinvoicer
Group=goinvoicer
WorkingDirectory=/srv/goinvoicer
EnvironmentFile=/srv/goinvoicer/goinvoicer.env
ExecStart=/srv/goinvoicer/goinvoicer
Restart=always
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/srv/goinvoicer

[Install]
WantedBy=multi-user.target
